# base-workspace/resources/etc/supervisor/conf.d/xdummy.conf
[program:xdummy]
command=/usr/lib/xorg/Xorg -noreset -nolisten tcp +extension GLX +extension RANDR +extension RENDER -logfile /tmp/xdummy.log :0.0
redirect_stderr=true 
stdout_logfile=/tmp/supervisor/%(program_name)s.log ; log logs into file
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected) 
startretries=5   ; max
autostart=false